# systemd/bt-watch.timer
[Unit]
Description=Start bt-watch on boot

[Timer]
OnBootSec=30
Unit=bt-watch.service

[Install]
WantedBy=timers.target
